# crontab
*	*	*	*	*	/var/www/html/vendor/bin/drush minter_trader:runTest --root=/var/www/html/web
